# testing/files/poweroff.timer
[Unit]
Description=Delayed poweroff

[Timer]
OnActiveSec=1
Unit=poweroff.target
